# scripts/inoxtags.ini
[program:inoxtags] 
command=/home/oromerob/webapps/inoxtags/scripts/start_inoxtags 
directory=/home/oromerob/webapps/inoxtags 
user=oromerob 
autostart=true 
autorestart=true 
redirect_stderr=True
stdout_logfile=/home/oromerob/logs/user/inoxtags.log
stderr_logfile=/home/oromerob/logs/user/inoxtags.log
